# files/etc/init.d/enable-rngd
#!/bin/sh /etc/rc.common

START=23

start() {
	uci set system.@rngd[0].enabled='1'
	uci set system.@rngd[0].device='/dev/urandom'
	uci set system.@rngd[0].fill_watermark='4000'
	uci commit
}
